# src/service/data/deepin-boot-maker.service
[Unit]
Description=Deepin Boot Maker
Wants=dbus.socket
After=dbus.socket

[Service]
Type=dbus
BusName=com.deepin.bootmaker
User=root
ExecStart=/usr/lib/deepin-daemon/deepin-boot-maker-service
StandardOutput=syslog
# Needs CAP_SYS_ADMIN umount u-disk.
CapabilityBoundingSet=~CAP_SYS_BPF CAP_NET_ADMIN
MemoryLimit=10G
IOWeight=200
ProtectSystem=full
ProtectHome=yes
InaccessiblePaths=/home
ReadWritePaths=-/dev/
ReadWritePaths=-/usr/lib/syslinux/mbr/

[Install]
WantedBy=multi-user.target
